# etc/crontab
# m h dom mon dow  command
* * * * * /home/pi/raspi/temp.py
