# BackupPi/etc/crontab
# m h dom mon dow  command
@reboot /home/pi/runWeeklyBackupAfter22h.sh
